# ops/codex-control.service
[Unit]
Description=Codex Telegram Control Plane
After=network.target

[Service]
Type=simple
User=codexsvc
Group=codexsvc
WorkingDirectory=/srv/codex/app/control-plane
EnvironmentFile=/srv/codex/app/control-plane/.env
ExecStart=/usr/bin/node dist/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
